# config/user_a_ble.conf
[program:user_a_ble]
command=/home/pi/radio-pi/user-a-ble/usera_peripheral.py
directory=/home/pi/radio-pi/user-a-ble/
environment=HOME="/home/pi"
user=root
priority=34
autostart=true
autorestart=true