# config/user_b_ble.conf
[program:user_b_ble]
command=/home/pi/radio-pi/user-b-ble/userb_peripheral.py
directory=/home/pi/radio-pi/user-b-ble/
environment=HOME="/home/pi"
user=root
priority=35
autostart=true
autorestart=true